# debian/postinst
#!/bin/sh -e

action="$1"
oldversion="$2"

umask 022

# source debconf stuff
. /usr/share/debconf/confmodule
# source dbconfig-common stuff
. /usr/share/dbconfig-common/dpkg/postinst.mysql

mkdir -p /etc/grase/

dbc_first_version=3.0
dbc_generate_include="template:/etc/grase/radmin.conf"
dbc_generate_include_perms="660"
dbc_generate_include_owner="root:www-data"
dbc_generate_include_args="-o template_infile=/usr/share/grase/templates/radmin.conf -U"

dbc_go grase-www-portal $@



if [ "$action" != configure ]
  then
  exit 0
fi

fix_dir_perms(){
	chgrp www-data -R /usr/share/grase/www/ || true
	chmod ag+rw -R /usr/share/grase/www/images/logo.png /usr/share/grase/www/radmin/templates_c /usr/share/grase/www/uam/templates_c || true
	chmod a+rw -R /usr/share/grase/www/public/* || true

}

copy_favicon(){
    cp -n /usr/share/grase/www/favicon.ico /var/www/favicon.ico || true
}

link_logo(){
    #if [[ ! -a /usr/share/grase/www/images/logo.png ]]
    #then
    
    # Don't worry about tests, will fail silently
    ln -s graselogo.png /usr/share/grase/www/images/logo.png > /dev/null || true
    #fi

}
#DEBHELPER#

link_logo

fix_dir_perms

copy_favicon

a2enmod rewrite || true

if which invoke-rc.d >/dev/null 2>&1; then
       invoke-rc.d apache2 restart 3>&- || true
else
       /etc/init.d/apache2 restart 3>&- || true
fi


exit 0


